# roles/privoxy/files/privoxy-blocklist.timer
[Unit]
Description=Monthly get new blocklists

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=multi-user.target
