# roles/privoxy/files/privoxy-blocklist.service
[Unit]
Description=Update privoxy blocklist

[Service]
Restart=no
ExecStart=/usr/local/bin/privoxy-blocklist.sh
